# deploy/launchd/co.poke.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent for the Beeper -> Poke bridge (macOS).

  Install:
    cp deploy/launchd/co.poke.bridge.plist ~/Library/LaunchAgents/
    # edit the WorkingDirectory and ProgramArguments paths below to match
    launchctl load ~/Library/LaunchAgents/co.poke.bridge.plist

  KeepAlive restarts the bridge if it exits or is killed. Beeper Desktop must be
  running for the local API to be reachable.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>co.poke.bridge</string>
    <key>WorkingDirectory</key>
    <string>/Users/YOU/poke-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOU/poke-bridge/.venv/bin/python</string>
        <string>bridge.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/YOU/poke-bridge/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/poke-bridge/launchd.err.log</string>
</dict>
</plist>
